# cloudflare-ddns-0.2.0/debian/cloudflare-ddns.timer
[Unit]
Description=Run cloudflare-ddns every hour

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
